# installation/defaults/desktop/config/systemd/user/ae-refresh-walker.path
[Unit]
Description=Watch Walker defaults for changes

[Path]
PathChanged=%h/.config/archenemy/installation/defaults/desktop/config/walker
PathModified=%h/.config/archenemy/installation/defaults/desktop/config/walker
Unit=ae-refresh-walker.service

[Install]
WantedBy=default.target
