# installation/defaults/desktop/config/systemd/user/ae-refresh-hyprland.path
[Unit]
Description=Watch Hyprland defaults for changes

[Path]
PathChanged=%h/.config/archenemy/installation/defaults/desktop/config/hypr
PathModified=%h/.config/archenemy/installation/defaults/desktop/config/hypr
Unit=ae-refresh-hyprland.service

[Install]
WantedBy=default.target
