# installation/defaults/system/power/systemd/user/battery-monitor.service
[Unit]
Description=Battery level watcher
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash -lc '
  battery=$(upower -e | grep -m1 BAT || true)
  if [[ -z "$battery" ]]; then
    exit 0
  fi
  state=$(upower -i "$battery" | awk "/state/ {print \$2}")
  level=$(upower -i "$battery" | awk "/percentage/ {gsub(\"%\", \"\", \$2); print int(\$2)}")
  if [[ "$state" == "discharging" && "$level" -le 20 ]]; then
    notify-send "Battery Low" "Charge at ${level}%%"
  fi
'
